# systemd/dot-config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/home/vandal/.local/bin/mailsync
